# deploy/cafeteria.service
[Unit]
Description=Cafeteria service

[Service]
ExecStart=/home/potados/cafeteria-server/run
WorkingDirectory=/home/potados/cafeteria-server
Restart=always
User=potados

[Install]
WantedBy=multi-user.target
